# stromzaehler_daily.timer
[Unit]
Description=Täglich stromzähler_daily.sh ausführen

[Timer]
OnCalendar=00:01

[Install]
WantedBy=timers.target
