# stromzaehler_daily.timer
[Unit]
Description=Execute stromzähler_daily.sh daily

[Timer]
OnCalendar=00:01

[Install]
WantedBy=timers.target
